# vagrant/systemd/elderberrypi.timer
[Unit]
Description=ElderberryPi timer

[Timer]
OnBootSec=10s
OnUnitInactiveSec=30min

[Install]
WantedBy=timers.target
